# trunk/packages/deb/prerm.in
#!/bin/bash

update-rc.d autopoweroff remove
rm -f @autopoweroff_sbindir@/autopoweroff
