# packages/deb/postinst.in
#!/bin/bash

# For nice examples on how to write this scripts, see scripts available
# under:  /var/lib/dpkg/info

#echo "args:  $*"
#echo "$0 called."
/usr/sbin/autopoweroff-upgrade
update-rc.d autopoweroff defaults >/dev/null
exit 0
